# templates/etc/cron.d/core_api_sync.in
*/5 * * * * www-data ${buildout:directory}/bin/ebay/manage ${buildout:directory}/src/inventorum.ebay/${config:environment}.ini runscript core_api_sync -v3
