# templates/etc/supervisor/staging.conf.in
[unix_http_server]
file=${os:root-prefix}/var/run/com.inventorum.ebay/supervisor.sock

[supervisord]
logfile = ${os:root-prefix}/var/log/com.inventorum.ebay/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = ${os:root-prefix}/var/run/com.inventorum.ebay/celery.supervisord.pid
nodaemon = false
# http://supervisord.org/configuration.html
environment = LD_LIBRARY_PATH="${buildout:directory}/parts/postgresql/lib"

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///${os:root-prefix}/var/run/com.inventorum.ebay/supervisor.sock


[group:celery]
programs=celery_worker_celery, celery_worker_fetching, celery_worker_updating, celery_worker_publishing, celery_worker_syncing, celery_beat, flower

[program:celery_worker_celery]
command = ${buildout:directory}/bin/ebay/celery src/inventorum.ebay/staging.ini -A inventorum.ebay.celery worker -Q celery -l info
process_name = com.inventorum.ebay.celery_worker_celery
directory = ${buildout:directory}
priority = 10
redirect_stderr = true
stdout_logfile = ${os:root-prefix}/var/log/com.inventorum.ebay/celery_worker.celery.log
autostart=true
autorestart=true

[program:celery_worker_fetching]
command = ${buildout:directory}/bin/ebay/celery src/inventorum.ebay/staging.ini -A inventorum.ebay.celery worker -Q fetching -l info
process_name = com.inventorum.ebay.celery_worker_fetching
directory = ${buildout:directory}
priority = 10
redirect_stderr = true
stdout_logfile = ${os:root-prefix}/var/log/com.inventorum.ebay/celery_worker.fetching.log
autostart=true
autorestart=true

[program:celery_worker_updating]
command = ${buildout:directory}/bin/ebay/celery src/inventorum.ebay/staging.ini -A inventorum.ebay.celery worker -Q updating -l info
process_name = com.inventorum.ebay.celery_worker_updating
directory = ${buildout:directory}
priority = 10
redirect_stderr = true
stdout_logfile = ${os:root-prefix}/var/log/com.inventorum.ebay/celery_worker.updating.log
autostart=true
autorestart=true

[program:celery_worker_publishing]
command = ${buildout:directory}/bin/ebay/celery src/inventorum.ebay/staging.ini -A inventorum.ebay.celery worker -Q publishing -l info
process_name = com.inventorum.ebay.celery_worker_publishing
directory = ${buildout:directory}
priority = 10
redirect_stderr = true
stdout_logfile = ${os:root-prefix}/var/log/com.inventorum.ebay/celery_worker.publishing.log
autostart=true
autorestart=true

[program:celery_worker_syncing]
command = ${buildout:directory}/bin/ebay/celery src/inventorum.ebay/staging.ini -A inventorum.ebay.celery worker -Q syncing -l info
process_name = com.inventorum.ebay.celery_worker_syncing
directory = ${buildout:directory}
priority = 10
redirect_stderr = true
stdout_logfile = ${os:root-prefix}/var/log/com.inventorum.ebay/celery_worker.syncing.log
autostart=true
autorestart=true

[program:celery_beat]
command = ${buildout:directory}/bin/ebay/celery src/inventorum.ebay/staging.ini -A inventorum.ebay.celery beat -l debug
process_name = celery_beat
directory = ${buildout:directory}
priority = 10
redirect_stderr = true
stdout_logfile = ${os:root-prefix}/var/log/com.inventorum.ebay/celery_beat.staging.log
autostart=true
autorestart=true

[program:flower]
command = ${buildout:directory}/bin/ebay/celery src/inventorum.ebay/staging.ini flower -A inventorum.ebay.celery --address=127.0.0.1 --port=5555
process_name = flower
directory = ${buildout:directory}
priority = 10
redirect_stderr = true
stdout_logfile = ${os:root-prefix}/var/log/com.inventorum.ebay/flower.staging.log
autostart=true
autorestart=true
